# xboxdrvDaemon.service
[Unit]
Description=yet another xboxdrv daemon

[Service]
ExecStart=/bin/bash /usr/local/bin/xboxdrvDaemon

[Install]
WantedBy=multi-user.target